# data_crawl/crontab
* * * * * sh -c "sleep 10s ;python3 /app/test.py" >>/tmp/out.log 2>/tmp/err.log
